# deploy/monger-frontend.service
[Unit]
Description=Monger Frontend (SvelteKit)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/monger-frontend
ExecStart=/root/.bun/bin/bun run build/index.js
Restart=on-failure
RestartSec=5
Environment=PORT=3000
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
