# deploy/monger-frontend.service
[Unit]
Description=Monger Frontend (SvelteKit + Node)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/monger-frontend/build
ExecStart=/usr/bin/node ./index.js
Restart=on-failure
RestartSec=5
Environment=PORT=3000
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
